# debian/spreedbox-authd.postinst
#!/bin/sh

set -e

case "$1" in
    configure)
        /usr/share/spreedbox-authd/spreedbox-authd-genkey-helper "$1" "$DEB_VERSION"

        if [ -e "/etc/init.d/nginx" -a -e "/var/run/nginx.pid" ]; then
            invoke-rc.d nginx reload || true
        fi
        ;;
    *)
        exit 0
    ;;
esac

#DEBHELPER#
